# raspberrypi/systemd/smartbp-ble.service
[Unit]
Description=SmartBP BLE Bridge
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/smartbp
ExecStart=/usr/bin/python3 /opt/smartbp/ble_bridge.py
Restart=always

[Install]
WantedBy=multi-user.target
